# ArVideoSubServer.conf

description "ArVideoSubServer"

console none
respawn
respawn limit 5 20

start on (local-filesystems and net-device-up IFACE!=lo)
stop on runlevel [016]

setgid root
setuid root

script
  /usr/local/Aria/ArVideo/examples/videoSubServer  >/var/log/ArVideoSubServer.log 2>&1
end script

