# roscore.conf

# Start roscore (ROS Hydro) when eth1 interface is available. IP address of eth1 is used in ROS_IP URI so other computers attached to eth1 can use this roscore. 

description "roscore (ROS hydro)"

console none
respawn
respawn limit 5 20

start on (local-filesystems and net-device-up IFACE!=lo)
stop on runlevel [016]

env ROS_LOG_DIR=/var/log/ros
#env ROS_IP=XXX

setgid ros
setuid ros

script
  cd /var/lib/ros
  . /opt/ros/hydro/setup.sh
  roscore
end script

